# supervisord.conf
[supervisord]

[program:uwsgi]
command=uwsgi -i uwsgi_crux.ini
user=ec2-user
